# flatpak/.config/systemd/user/unattended-flatpak-update.timer
[Unit]
Description=Update flatpak user packages

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
