# flatpak/.config/systemd/user/unattended-flatpak-update.service
[Unit]
Description=Unattended flatpak update

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak --user update --assumeyes --noninteractive
