# common/.config/systemd/user/spotifyd.service
[Unit]
Description=Spotify Daemon

[Service]
ExecStart=/usr/bin/spotifyd --no-daemon

[Install]
WantedBy=default.target
