# sway/.config/systemd/user/gsettings-preferences.service
[Unit]
Description=gsettings preferences

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gsettings-preferences.sh

[Install]
WantedBy=sway.service
