# sway/.config/systemd/user/sway.service
[Unit]
Description=sway - SirCmpwn's Wayland window managers
Documentation=man:sway(5)

[Service]
Type=simple
ExecStart=/usr/bin/sway

