# nvim/.config/systemd/user/unattended-nvim-update.service
[Unit]
Description=Unattended nvim plugins update

[Service]
Type=oneshot
ExecStart=/usr/bin/nvim --headless +"PlugUpgrade" +"PlugUpdate!" +"qa!"
