# nvim/.config/systemd/user/unattended-nvim-update.timer
[Unit]
Description=Update neovim plugins

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
